# systemd/gitops-pull.timer
[Unit]
Description=GitOps pull timer — check for infrastructure changes every 1 minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
